# infra/ansible/playbooks/templates/services/node_exporter.service.j2
[Unit]
Description=Node exporter
After=network.target

[Service]
Type=simple
User={{ ansible_user }}
WorkingDirectory=/home/{{ ansible_user }}/
ExecStart=node_exporter --collector.systemd --collector.processes
Restart=always
RestartSec=1
StartLimitBurst=100
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
